# src/etc/systemd/eardbd.service.in
[Unit]
Description=EARDBD - Energy Aware Runtime database cache daemon
ConditionPathExists=@localstatedir@
#Wants=req1.service req2.service
#After=req1.service req2.service


[Service]
Type=simple
Environment="EAR_ETC=@sysconfdir@"
ExecStart=@prefix@/sbin/eardbd
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
#Restart=on-failure
KillMode=process
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target

