# etc/systemd/edcmon.service.in
[Unit]
Description=EDCMON - EAR Data Center Monitor service
#ConditionPathExists=@localstatedir@
#Wants=req1.service req2.service
#After=req1.service req2.service


[Service]
Type=simple
Environment="EAR_ETC=@sysconfdir@"
Environment="EAR_INSTALL_PATH=@prefix@"
# nodesensors_report period is not used but the API requires it. The nodesensor_log is the argument showing which report plugin to use
ExecStart=@prefix@/sbin/edcmon --monitor=1000 --relax=1000 --plugins=nodesensors.so:30000+nodesensors_report.so:30000:nodesensor_log+nodesensors_alerts.so:30000:log --verbose
ExecStop=/bin/kill -TERM $MAINPID
KillMode=process
Restart=on-failure
#RestartSec=1s
LimitMEMLOCK=infinity
LimitSTACK=infinity

[Install]
WantedBy=multi-user.target

